# youtube-worker.service
# ============================================================
# YouTube Transcription Worker — systemd 用户服务
# 消费 youtube_fetching + youtube_transcription 两个队列
#
# 安装：
#   mkdir -p ~/.config/systemd/user
#   cp youtube-worker.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now youtube-worker.service
# ============================================================

[Unit]
Description=YouTube Transcription Worker (youtube_fetching + youtube_transcription)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/local_virtual_service
EnvironmentFile=%h/local_virtual_service/.env
Environment=PYTHONPATH=%h/local_virtual_service

ExecStart=%h/local_virtual_service/venv/bin/celery -A worker.celery_app worker \
    --queues=youtube_fetching,youtube_transcription \
    --concurrency=2 \
    --max-tasks-per-child=20 \
    --hostname=youtube-transcription-worker@%%h \
    --loglevel=info \
    --logfile=%h/local_virtual_service/logs/worker.log \
    --pidfile=%h/local_virtual_service/logs/worker.pid

Restart=always
RestartSec=15

[Install]
WantedBy=default.target
